# context/etc/supervisor/supervisord.conf
[supervisord]
nodaemon=true

[program:express]
command=node /www-data/index.js

[program:openarena]
command=/opt/start_openarena.sh
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0